# deploy/systemd/polymarket-sentinel.service
[Unit]
Description=Polymarket Sentinel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/polymarket-sentinel
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/polymarket-sentinel/.venv/bin/python -m sentinel
Restart=always
RestartSec=10
TimeoutStopSec=20
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
